# bin/rc.d/alignak-webui
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: alignak_webui
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="alignak_webui"
rcvar="alignak_webui_enable"

alignak_webui_configfile="/usr/local/etc/alignak-webui/uwsgi.ini"
command="/usr/local/bin/uwsgi"
command_interpreter="/usr/local/bin/python2.7"
pidfile="/var/run/alignak_webui.pid"
command_args="--ini ${alignak_webui_configfile} -d /dev/null --pidfile ${pidfile}"

required_files="${alignak_webui_configfile}"

load_rc_config "${name}"

[ -z "${alignak_webui_enable}" ] && alignak_webui_enable="NO"

run_rc_command "$1"
