# bin/systemd/alignak-webui.service
[Unit]
Description=uWSGI instance to serve Alignak WebUI
After=network.target

[Service]
# Environment variables
Environment=ALIGNAK_WEBUI_CONFIGURATION_FILE=/usr/local/etc/alignak-webui/settings.cfg
Environment=ALIGNAK_WEBUI_BACKEND=http://127.0.0.1:5000
Environment=ALIGNAK_WEBUI_WS=http://127.0.0.1:8888
Environment=ALIGNAK_WEBUI_LOG=/usr/local/var/log/alignak-webui
Environment=ALIGNAK_WEBUI_LOGGER_FILE=/usr/local/etc/alignak-webui/logging.json
# Execution
ExecStartPre=/bin/mkdir -p /usr/local/var/log/alignak-webui
ExecStartPre=/bin/chown alignak:alignak /usr/local/var/log/alignak-webui
ExecStart=/usr/bin/uwsgi --ini /usr/local/etc/alignak-webui/uwsgi.ini
KillSignal=SIGINT
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
